# projects/kazimierz.akn/src/infrastructure/ansible/roles/pangolin/templates/geoip-update.service.j2
[Unit]
Description=Update Pangolin GeoIP database
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-pangolin-geoip
User=root
